# packaging/systemd/ya-router.service
[Unit]
Description=ya-router OpenAI-compatible routing daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=ya-router
Group=ya-router
WorkingDirectory=/var/lib/ya-router
EnvironmentFile=-/etc/ya-router/ya-router.env
Environment=YA_ROUTER_CONFIG_DIR=/var/lib/ya-router
Environment=YA_ROUTER_CONFIG_PATH=/var/lib/ya-router/config.json
Environment=YA_ROUTER_SECRETS_PATH=/var/lib/ya-router/secrets.json
Environment=YA_ROUTER_OPERATIONS_PATH=/var/lib/ya-router/operations.json
Environment=YA_ROUTER_COOLDOWNS_PATH=/var/lib/ya-router/cooldowns.json
Environment=YA_ROUTER_CONTROL_SOCKET=/run/ya-router/control.sock
Environment=YA_ROUTER_LISTEN_ADDRESS=127.0.0.1
ExecStart=/usr/local/bin/ya-routerd run --config-migrate=none
Restart=on-failure
RestartSec=5s
TimeoutStopSec=45s
KillSignal=SIGTERM
StateDirectory=ya-router
StateDirectoryMode=0700
RuntimeDirectory=ya-router
RuntimeDirectoryMode=0700
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
